# ssy-easemob-arvon/201606/ops-repo/roles/common/ntp/templates/etc/cron.d/ntp
*/5 * * * * root /usr/sbin/ntpdate {{ntpserver}}
